# post/base/etc/systemd/system/update.timer
[Unit]
Description=Eyerise auto system update hourly

[Timer]
OnCalendar=hourly
Persistent=true
Unit=update.service

[Install]
WantedBy=timers.target