# src/base/etc/systemd/system/update.service
[Unit]
Description=Creamie auto system update

[Service]
Type=oneshot
ExecStart=/usr/bin/pacman --sync --refresh --sysupgrade --noconfirm